# debian/tofurkey.service
# Example systemd unit file for tofurkey
[Unit]
Description=tofurkey
Documentation=man:tofurkey
After=local-fs.target
After=network.target
After=time-sync.target

[Install]
WantedBy=multi-user.target

[Service]
Type=notify
Restart=on-failure
ExecStart=/usr/sbin/tofurkey
